# supervisord.conf
[supervisord]
nodaemon=true

[program:process1]
command=java -jar /opt/app.jar

[program:process2]
command=tor